# root/etc/s6-overlay/s6-rc.d/socat-imap/run
#!/command/with-contenv sh

exec s6-setuidgid socat TCP-LISTEN:143,fork TCP:127.0.0.1:1143
