# root/etc/s6-overlay/s6-rc.d/proton-bridge/run
#!/command/with-contenv sh

cd /opt/protonmail

exec cat /opt/protonmail/etc/protontty | su-exec appuser /opt/protonmail/bin/proton-bridge --noninteractive
