# root/etc/s6-overlay/s6-rc.d/proton-bridge/run
#!/command/with-contenv sh

exec cat /opt/protonmail/etc/protontty | s6-setuidgid appuser /opt/protonmail/bin/proton-bridge --noninteractive
