# root/etc/s6-overlay/s6-rc.d/socat-smtp/run
#!/command/with-contenv sh

exec socat TCP-LISTEN:25,fork TCP:127.0.0.1:1025
